# deploy/systemd/luckyme-push-alerts.service
[Unit]
Description=LuckyMe round push alerts
After=network-online.target luckyme-api.service
Wants=network-online.target

[Service]
Type=oneshot
User=luckyme
Group=luckyme
WorkingDirectory=/opt/luckyme
EnvironmentFile=/etc/luckyme/luckyme-api.env
Environment=CONFIRM_MAINNET_PUSH_ALERTS=true
Environment=LUCKYME_PUSH_SEND=true
Environment=NPM_CONFIG_UPDATE_NOTIFIER=false
ExecStart=/usr/bin/npm run push:round-alerts
NoNewPrivileges=true
PrivateTmp=true
