# deploy/systemd/luckyme-promotional-pools-keeper.service
[Unit]
Description=LuckyMe promotional pools MagicBlock keeper
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=luckyme
Group=luckyme
WorkingDirectory=/opt/luckyme
EnvironmentFile=/etc/luckyme/admin-promotions.env
Environment=LUCKYME_PROMOTIONS_KEEPER_MAX_ACTIONS=4
ExecStart=/usr/bin/node /opt/luckyme/scripts/promotional-pools-keeper.mjs
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/luckyme-promotions
