# deploy/systemd/luckyme-operations-monitor.timer
[Unit]
Description=Run LuckyMe operational monitor every minute

[Timer]
OnBootSec=120
OnUnitActiveSec=60
AccuracySec=10
Persistent=true
Unit=luckyme-operations-monitor.service

[Install]
WantedBy=timers.target
