# deploy/systemd/luckyme-settlement-keeper.timer
[Unit]
Description=Run LuckyMe settlement keeper every minute

[Timer]
OnActiveSec=5
OnUnitInactiveSec=60
AccuracySec=5
Persistent=true
Unit=luckyme-settlement-keeper.service

[Install]
WantedBy=timers.target
