# deploy/systemd/luckyme-push-alerts.timer
[Unit]
Description=Run LuckyMe round push alerts every minute

[Timer]
OnBootSec=90
OnUnitActiveSec=60
AccuracySec=10
Persistent=true
Unit=luckyme-push-alerts.service

[Install]
WantedBy=timers.target
