# deploy/systemd/luckyme-seeker-referral.service
[Unit]
Description=LuckyMe Seeker referral API
After=network-online.target luckyme-api.service
Wants=network-online.target

[Service]
Type=simple
User=luckyme
Group=luckyme
WorkingDirectory=/opt/luckyme-seeker-referral
EnvironmentFile=/etc/luckyme/luckyme-api.env
Environment=NODE_ENV=production
Environment=REFERRAL_REQUIRE_HTTPS=true
Environment=REFERRAL_TRUST_PROXY=true
Environment=REFERRAL_TEST_MODE=false
Environment=APP_ANALYTICS_ENABLED=true
Environment=SEEKER_PASS_PROMOTION_ENABLED=true
Environment=SEEKER_REFERRAL_HOST=127.0.0.1
Environment=SEEKER_REFERRAL_PORT=8790
Environment=SEEKER_REFERRAL_DB_PATH=/var/lib/luckyme-seeker-referral/referral.sqlite
Environment=LUCKYME_SETTLEMENT_ARCHIVE_PATH=/var/lib/luckyme/settlements.jsonl
ExecStart=/usr/bin/node backend/src/seeker-referral-server.mjs
Restart=always
RestartSec=5
StateDirectory=luckyme-seeker-referral
StateDirectoryMode=0750
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/luckyme-seeker-referral

[Install]
WantedBy=multi-user.target
